# deploy/systemd/crypto-bot-watchdog.service
[Unit]
Description=Crypto Trading Machine watchdog (restarts crypto-bot if the autopilot scheduler stalls)
After=network.target crypto-bot.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/crypto-trading-machine
ExecStart=/opt/crypto-trading-machine/scripts/watchdog.sh
